# docker/pythontask.timer
[Timer]
OnBootSec=5min
OnUnitActiveSec=1d

